# dt_image_search/deploy/systemd/imagesearch.service
# /etc/systemd/system/myapp.service
[Unit]
Description=Image-Search Docker Compose Service
Requires=docker.service
After=docker.service

[Service]
WorkingDirectory=/path/to/your/project
ExecStart=/usr/local/bin/docker-compose up -d
ExecStop=/usr/local/bin/docker-compose down
Restart=always
RestartSec=10s
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target