# web/instant-share/deploy/instant-share-relay.service
[Unit]
Description=Instant-Share WebSocket Relay
After=network.target

[Service]
ExecStart=/opt/instant-share-relay/start.sh
Restart=always
Environment=RELAY_PORT=8787

[Install]
WantedBy=multi-user.target